# storage-server/crontab
cat /etc/crontab

# *  *  *  *  * user-name  command to be executed
5 18 * * * root /usr/bin/bash /etc/scripts/ldap_backup.sh >> /var/log/cron 2>&1
