# ldap-server/crontab
cat /etc/crontab

# *  *  *  *  * user-name  command to be executed

0 11 * * * root /usr/bin/bash /root/backup/bkp_commands.sh
0 18 * * * root /usr/bin/bash /root/backup/bkp_commands.sh
